# build/install/docker/config/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:www]
environment=NODE_ENV=production,NODE_TLS_REJECT_UNAUTHORIZED=0,NODE_OPTIONS="--dns-result-order=ipv4first"
directory=/var/www/onlyoffice/controlpanel/www/
command=node server.js
autostart=true
autorestart=true
